# supervisor.conf
[supervisord]
nodaemon=true

[program:macudev]
command=node /app/server/index.js
stdout_logfile_maxbytes = 0
redirect_stderr=true
stdout_logfile=/dev/sdtout
